# workflow-demo/shiny-container/supervisord.conf
[supervisord]
nodaemon=true

[program:app.py]
command=python /app/app.py

[program:shiny-server]
command=/usr/bin/shiny-server.sh